# authelia/rootfs/etc/s6-overlay/s6-rc.d/init-authelia/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Authelia
# Configures Authelia before running
# ==============================================================================

bashio::log.debug "Initializing Authelia environment variables..."

# Set environment variables
export AUTHELIA_SERVER_ADDRESS="tcp://:9091"
export AUTHELIA_STORAGE_LOCAL_PATH="/data/db.sqlite3"

# Log the values of the environment variables
bashio::log.debug "AUTHELIA_SERVER_ADDRESS set to: ${AUTHELIA_SERVER_ADDRESS}"
bashio::log.debug "AUTHELIA_STORAGE_LOCAL_PATH set to: ${AUTHELIA_STORAGE_LOCAL_PATH}"

bashio::log.debug "Starting Authelia configuration check..."

# Define paths
CONFIG_FILE="/config/config.yml"
TEMPLATE_FILE="/etc/authelia/config.template.yml"

# Ensure configuration exists
if ! bashio::fs.file_exists "$CONFIG_FILE"; then
    bashio::log.warning "Configuration file '$CONFIG_FILE' not found!"
    
    if bashio::fs.file_exists "$TEMPLATE_FILE"; then
        bashio::log.info "Copying template configuration from '$TEMPLATE_FILE' to '$CONFIG_FILE'..."
        cp "$TEMPLATE_FILE" "$CONFIG_FILE"
        
        if [[ $? -eq 0 ]]; then
            bashio::log.info "Configuration file successfully created at '$CONFIG_FILE'."
            bashio::log.warning "Please edit '$CONFIG_FILE' with your settings before restarting."
            exit 1  # Exit to indicate the first-time setup
        else
            bashio::log.error "Failed to copy the template configuration file!"
            exit 1  # Exit due to an error
        fi
    else
        bashio::log.error "Template configuration file '$TEMPLATE_FILE' not found!"
        exit 1  # Exit due to missing template
    fi
else
    bashio::log.debug "Existing configuration found at '$CONFIG_FILE'. Continuing startup..."
fi